# 3.10/rootfs/etc/supervisor.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
redirect_stderr=true
autostart=true
stdout_syslog=true

