# 3.3/supervisor/supervisor.d/crond.conf
[program:crond]
command=/usr/sbin/crond -f
